# deb/DEBIAN/postrm
#!/bin/sh

update-rc.d -f logwall remove
